# etc/supervisord.d/sensu-server.ini
[program:sensu-server]
command=/opt/sensu/embedded/bin/ruby /opt/sensu/bin/sensu-server -c /etc/sensu/config.json -d /etc/sensu/conf.d -e /etc/sensu/extensions -p /var/run/sensu/sensu-server.pid -l /var/log/sensu/sensu-server.log -L info
user=sensu
autostart=true
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
