# etc/supervisord.d/sensu-client.ini
[program:sensu-client]
command=/opt/qnib/sensu/client/bin/start.sh
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=1
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
